# deploy/finance-toolkit.service
[Unit]
Description=Finance Team Toolkit (FastAPI / Gunicorn+Uvicorn)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/finance-toolkit
Environment="PATH=/opt/finance-toolkit/.venv/bin"
# --forwarded-allow-ips makes the uvicorn worker trust X-Forwarded-* from
# Nginx (TLS termination), so scheme/secure cookies work. NB: --proxy-headers
# is a uvicorn-CLI flag, NOT a gunicorn one — don't add it here.
ExecStart=/opt/finance-toolkit/.venv/bin/gunicorn app.main:app \
    -k uvicorn.workers.UvicornWorker \
    -w 2 -b 127.0.0.1:8801 \
    --forwarded-allow-ips=127.0.0.1 \
    --access-logfile - --error-logfile -
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
